# usr/lib/opkg/info/cc-client.list
/usr/sbin/C2DISCONNECT
/usr/sbin/C2CONNECT
/usr/sbin/C2EXFIL
/etc/init.d/cc-client
/usr/sbin/cc-client
