# usr/lib/opkg/info/pineap.list
/usr/bin/pineap
/etc/init.d/pineapd
/usr/sbin/pineapd
