# usr/lib/opkg/info/resetssids.list
/etc/init.d/resetssids
/usr/sbin/resetssids
