# usr/lib/opkg/info/rtl-sdr.list
/usr/bin/rtl_sdr
/usr/bin/rtl_tcp
/usr/bin/rtl_fm
/etc/config/rtl_tcp
/usr/bin/rtl_test
/usr/bin/rtl_power
/usr/bin/rtl_adsb
/etc/init.d/rtl_tcp
/usr/bin/rtl_eeprom
